# 7.0/nginx/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:php-fpm]
command=php-fpm7.0

[program:nginx]
command=nginx -g "daemon off;"
